# root/etc/s6-overlay/s6-rc.d/laravel-task-scheduler/run
#!/bin/sh

cd /var/www/html && php artisan schedule:work
