# root/etc/s6-overlay/s6-rc.d/laravel-queue-worker/run
#!/bin/sh

cd /var/www/html && php artisan queue:work --sleep=5 --tries=3
